# bin/debian-shutdown
#! /bin/sh
### BEGIN INIT INFO
# Provides:          heartbeat-shutdown
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Heartbeat initscript
# Description:       This file starts and stops the heartbeat daemon.
### END INIT INFO

SCRIPT_PATH=`/bin/readlink -f "$0"`
BASE_PATH=`/usr/bin/dirname "$SCRIPT_PATH"`/..

USER=root
GROUP=root

PATH=/sbin:/usr/sbin:/bin:/usr/bin

HEARTBEAT=$BASE_PATH/bin/heartbeat
PIDFILE=/tmp/heartbeat.pid

[ -x "$HEARTBEAT" ] || exit 0

. /lib/init/vars.sh
. /lib/lsb/init-functions

do_stop()
{
  echo "Running the heartbeat shutdown script."
  
  /bin/su -l $USER -c "$HEARTBEAT --config $BASE_PATH/config/shutdown.yml"
}

case "$1" in
  start)
  ;;
  stop)
    do_stop
  ;;
  restart|force-reload)
  ;;
  *)
    echo "Usage: heartbeat-shutdown {stop}"
    exit 3
  ;;
esac

:
